# wwa_devops.timer
[Unit]
Description=Run WWA Console tool every five minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target

