# wwa_devops.service
[Unit]
Description=WWA Console tool

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wwa_devops
User=root

